# v.conf
[program:v]
command=python3 /home/curator2/v/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/v/test.err.log
stdout_logfile=/home/curator2/v/test.out.log